# openmediavault/homecloud-0.68-arm64/etc/systemd/system/password-reset.service
[Unit]
Description=Password Reset Service
After=network.target docker.service openmediavault-engined.service oled-display-disk.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/password-reset
ExecStartPre=/etc/password-reset/start.sh
ExecStart=/usr/bin/docker compose -f /etc/password-reset/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /etc/password-reset/docker-compose.yml down
TimeoutStartSec=0

# Restart policy
Restart=on-failure
RestartSec=10s

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

